# deploy/systemd/autoone.service
[Unit]
Description=AutomationOne Stack (Docker Compose)
Documentation=https://linear.app/autoone/issue/AUT-336
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/robin/autoone
# Profile wie in .env (COMPOSE_PROFILES); bei Bedarf anpassen.
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.pi.yml --profile monitoring up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.pi.yml --profile monitoring stop
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
